# ryzen_temps.cron
# Ryzen CPU Temperature Logging
# Collect temperature data every minute
* * * * * root /root/ryzen_temps.sh >/dev/null 2>&1
